# fs-overlay/etc/inittab
::sysinit:/etc/init.d/rcS
tty1::respawn:-/bin/bash
ttyS0::respawn:-/bin/bash
::shutdown:/bin/sh /rescueOS/shutdown-cleanup 
